# etc/services/sensors.service
[Unit]
Description=Service that launches the sensors component.
After=network.target

[Service]
WorkingDirectory=/home/robocomp/robocomp/components/EBO_V2/components/sensors
ExecStart=/home/robocomp/software/pip_venv/bin/python src/sensors.py etc/config
Restart=always
RestartSec=5
User=robolab

[Install]
WantedBy=multi-user.target

